# other_modules/toms_modules/mysql/templates/profile.d/mysql.sh.erb
#!/bin/bash
export CDPATH=$CDPATH:<%= scope.function_extlookup(["common::mysql::cdpath"]) %>
