# other_modules/toms_modules/smartfox/templates/redbox.erb
#! /bin/sh
#
# This file Autogenerated for <%= fqdn %> by puppet. Hand edits considered harmful.
#
# $Id: redbox.tpl,v 1.7 2010/05/26 14:44:45 jjoy Exp $
#
# Redbox Init File
#
# chkconfig: - 55 45
# description:	The redbox java "daemon" does something or other for smartfox
# processname: java
# config: /etc/sysconfig/redbox
# pidfile: /usr/local/smartfox/SFS/Server/logs/redbox.pid
#

# Source function library.
. /etc/rc.d/init.d/functions

#  Set some vars

PROCESSNAME=redbox
APP_PATH=/usr/local/smartfox/SFS/Server/RedBox
JAVA_HOME=/usr/local/smartfox/SFS_REDBOX/jre
JAVA=$JAVA_HOME/bin/java
PATH=${PATH}:${APP_PATH}
REDBOX_HOME=/usr/local/smartfox/SFS/Server/RedBox
SMARTFOX_LOGS=/usr/local/smartfox/SFS/Server/logs
RUNDIR=${SMARTFOX_LOGS}

cd $REDBOX_HOME

start()
{
    if [ -e ${RUNDIR}/${PROCESSNAME}.pid ] && [ -d /proc/`cat ${RUNDIR}/${PROCESSNAME}.pid` ] ; then
        echo "Process Already Running."
    else
        echo -n $"Starting ${PROCESSNAME}: "
        $JAVA -Djava.security.manager -Djava.security.policy=conf/red5.policy -cp red5.jar:conf:$CLASSPATH org.red5.server.Standalone 2>&1 1> /var/log/redbox.log &
        sleep 10;
        RETVAL=$?
        echo $!
        echo $! > ${RUNDIR}/${PROCESSNAME}.pid
        [ $RETVAL -eq 0 ] && touch /var/lock/subsys/${PROCESSNAME}
    fi
}

stop()
{
    echo -n $"Stopping ${PROCESSNAME}: "
    killproc -p $RUNDIR/${PROCESSNAME}.pid ${PROCESSNAME}
    RETVAL=$?
    echo
    if [ $RETVAL -eq 0 ] ; then
	rm -f /var/lock/subsys/${PROCESSNAME}
	rm -f $RUNDIR/${PROCESSNAME}.pid
    fi
}

restart () {
    stop
    sleep 1
    start
}

case "$1" in
    start)
	start
	;;
    stop)
	stop
	;;
    status)
	status ${PROCESSNAME}
	;;
    restart|reload)
	restart
	;;
    condrestart)
	[ -f /var/lock/subsys/${PROCESSNAME} ] && restart || :
	;;
    *)
	echo $"Usage: $0 {start|stop|status|restart|reload|condrestart}"
	exit 1
esac

exit $?
